# services/nginx/ssl/service_ssl
#!/sbin/openrc-run

# Service Definition file compatible with open rc for ssl service

start() {
    ebegin "Starting SSL Service"
    start-stop-daemon -S -b -m -p /var/wp/ssl/ssl.pid -x /var/wp/ssl/exe main /var/wp/ssl
    eend $?
}

stop() {
    ebegin "Stoping SSL Service"
    start-stop-daemon -K -p /var/wp/ssl/ssl.pid -x /var/wp/ssl/exe main /var/wp/ssl
    eend $?
}
